# templates/default/iscsiadm.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:        iscsiadm
# Required-Start:  $network $remote_fs
# Required-Stop:   $network $remote_fs sendsigs open-iscsi
# Default-Start:   S
# Default-Stop:    0 1 6
# Short-Description: Start iscsi attachements
### END INIT INFO

# *****************************************************************
# Dynamic iscsiadm configuration  --  DO NOT EDIT THIS FILE BY HAND
# Generated by chef - <%= @date %>
# *****************************************************************
PATH=/sbin:/bin:/usr/sbin:/usr/bin

_attach() {
 [ ! -z "$1" ] || return
 local ip=$1 i j getit; shift
 for i in $(iscsiadm --mode discovery --type sendtargets --portal $ip| grep "$ip:"| cut -d' ' -f2); do

  getit=1; for j in $*; do
   echo "$i"| grep -wqs "$j" && getit=1 && break
   getit=0
  done; [ $getit -ne 1 ] && continue

  iscsiadm --mode node --targetname $i --portal $ip --login || break

 done
}

do_start() {
 <%@portal.each do |i|%>
  _attach <%="#{i['ip']}"%> <%i['targets'].each do |target|%><%="\"#{target}\" "%><%end if i['targets']-%> || return
 <%end-%>
}

do_stop()  {
 local i; for i in $(iscsiadm -m session| cut -d' ' -f3,4| sed -e 's/ /|/g'); do
  iscsiadm --mode node -u --portal $(echo $i| cut -d'|' -f1| sed -e 's/:.*$//') -T $(echo $i| cut -d'|' -f2)
 done
}

do_status() {
 iscsiadm -m session
}

case "$1" in
start)
 if ! do_status >/dev/null 2>&1; then
  do_start || (do_stop && echo "$(basename $0): cannot start..." >&2 && false)
 else
  echo "Already started..." >&2
 fi
 ;;
stop)
 do_stop ;;
status)
 do_status ;;
*)
 echo "Syntaxe: $(basename $0) [start|stop|status]..." && false
esac
